# roles/configure_jumpbox/templates/cf_dyndns.timer.j2
[Unit]
Description="Timer for cloudflare dns renew script"

[Timer]
Unit=cf_dyndns.service
Persistent=true
# Run every day at 4am
OnCalendar=*-*-* 04:20:00

[Install]
WantedBy=timers.target
