# roles/configure_allowlistapp/files/allowlistapp.service
[Unit]
Description=allowlistapp IP allowlist app for nginx
After=network.target

[Service]
User=nginx
WorkingDirectory=/opt/allowlistapp
ExecStart=/usr/local/bin/poetry run flask --app allowlist run
Restart=always

[Install]
WantedBy=multi-user.target
