# roles/configure_elk/files/certbot.service
[Unit]
Description="Timer for certbot service"
Requires=certbot.timer

[Service]
Type=simple
ExecStart=/opt/scripts/certbot_elastic.sh
