# roles/configure_jumpbox/files/duckdns_renew.service
[Unit]
Description="Service to renew duck dns"
Requires=duckdns_renew.timer

[Service]
Type=simple
ExecStart=/opt/duckdns/duck.sh
