# roles/configure_jumpbox/files/duckdns_renew.service
[Unit]
Description="Timer for certbot renew service"
Requires=duckdns_renew.timer

[Service]
Type=simple
ExecStart=/opt/duckdns/duck.sh
