# roles/configure_reverse_proxy/files/certbot_renew.service
[Unit]
Description="Timer for certbot renew service"
Requires=certbot_renew.timer

[Service]
Type=simple
ExecStart=/usr/bin/certbot renew
