# roles/configure_reverse_proxy/files/certbot_renew.service
[Unit]
Description="Timer for duck dns renew service"
Requires=duckdns_renew.timer

[Service]
Type=simple
ExecStart=/usr/bin/certbot renew
